# disk-status/disk-status.service
[Unit]
Description=status reporting of disks

[Service]
ExecStart=/opt/clang/disk-status.sh 
User=root
