# backupserver/samba-workaround.timer
[Unit]
Description=Workaround to start samba

[Timer]
OnBootSec=4

[Install] 
WantedBy=timers.target
