# piserver/backup-server.timer
[Unit]
Description=Startup backup server

[Timer]
OnCalendar=Wed *-*-* 18:00:00
Persistent=true

[Install] 
WantedBy=timers.target
