# piserver/backup-server.service
[Unit]
Description=Startup backup server

[Service]
ExecStart=/opt/clang/backup-server.sh
User=root
