# piserver/generate-wiki.service
[Unit]
Description=Newly create the wiki.

[Service]
ExecStart=/opt/clang/generate-wiki.sh
User=root
