# piserver/generate-wiki.timer
[Unit]
Description=Trigger the generation of the wiki daily.

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
