# piserver/generate-wiki.service
[Unit]
Description=Newly create the wiki.

[Service]
ExecStart=/home/christian/system-config/generate-wiki.sh
User=root
