# mainserver/backup-pull.service
[Unit]
Description=daily backup pull task

[Service]
ExecStart=/opt/clang/backup-pull.sh
User=clang
