# mainserver/backup-pull.service
[Unit]
Description=Sync pull backup of piserver

[Service]
ExecStart=/opt/clang/backup-pull.sh
User=clang
