# mainserver/btrfs-status.service
[Unit]
Description=status reporting of btrfs filesystem

[Service]
ExecStart=/opt/clang/btrfs-status.sh 
User=root
