# mainserver/btrfs-status.service
[Unit]
Description=Report status of filesystem.

[Service]
ExecStart=/opt/clang/btrfs-status.sh 
User=root
