# mainserver/samba-workaround.timer
[Unit]
Description=Workaround to start samba

[Timer]
OnBootSec=20

[Install] 
WantedBy=timers.target
