# mainserver/snapshots.timer
[Unit]
Description=Daily snapshot task

[Timer]
OnCalendar=daily
Persistent=true

[Install] 
WantedBy=timers.target
