# mainserver/snapshots.service
[Unit]
Description=Perform snapshoting tasks

[Service]
ExecStart=/opt/clang/snapshots.sh 
User=root
