# mainserver/synced-media.service
[Unit]
Description=daily synced media task

[Service]
ExecStart=/opt/clang/synced-media.sh
User=clang
