# lang-backupserver/backup-pull.service
[Unit]
Description=Sync pull backup of main-server

[Service]
ExecStart=/opt/clang/backup-pull.sh
User=clang
