# lang-backupserver/backup-pull.timer
[Unit]
Description=weekly backup pull task

[Timer]
OnCalendar=weekly
Persistent=true

[Install] 
WantedBy=timers.target
